# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=gunicorn app:app

[program:bot]
command=python3 bot.py
